# systemd/claude-longrun-heartbeat.timer
[Unit]
Description=Fire claude long-run heartbeat at usage-limit reset grid (+6 min)

[Timer]
OnCalendar=*-*-* 02,07,11,16,21:06:00
Persistent=true

[Install]
WantedBy=timers.target
