# systemd/claude-longrun-heartbeat.service
[Unit]
Description=Claude long-run pipeline recovery heartbeat

[Service]
Type=oneshot
ExecStart=/bin/bash __HOME__/.claude/hooks/long-run-os-heartbeat.sh
