# systemd/claude-longrun-heartbeat.timer
[Unit]
Description=Fire claude long-run heartbeat every 30 min (gated by nextResetEpoch in the state file)

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target
